# units/wait-sni-ready.service.in
[Unit]
Description="Wait until tray dbus socket appeared"

[Service]
Type=oneshot
ExecStart=@execdir@/wait-sni-ready
RemainAfterExit=yes